# lib/systemd/system/steambuntu-autorepair.service
[Unit]
Description=SteamBuntu Autorepair

[Service]
ExecStart=/usr/bin/steambuntu-autorepair
Type=oneshot

